# frpc/amd64/etc/supervisor.d/docker-gen.conf
[program:docker-gen]
command=docker-gen -notify "supervisorctl restart frpc" -watch /etc/frp/frpc.tmpl /etc/frp/frpc.ini

