# frpc/amd64/etc/supervisor.d/frpc.conf
[program:frpc]
command=frpc -c /etc/frp/frpc.ini

